# supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
command=/usr/lib/postgresql/9.4/bin/pg_ctl start -D /etc/postgresql/9.4/main/; psql -c "CREATE ROLE root SUPERUSER LOGIN;"
